# airpi/scripts/profile.d/airpi_motd.sh
#!/bin/ash
[ "$TERM" = "dumb" ] && exit 0
/etc/bannerwrt.sh
/etc/info.sh
